# howdyd/systemd/howdyd.service
[Unit]
Description=Howdy daemon

[Service]
Type=dbus
BusName=nl.boltgolt.Howdy
ExecStart=/usr/libexec/howdyd
